# docker/php/crontab
* * * * * root /usr/local/bin/php /var/www/catalog/artisan schedule:run
